# deployment/systemd/esab-recommender.target
[Unit]
Description=ESAB Recommender Application (Frontend + Backend)
After=network.target postgresql.service neo4j.service
Wants=postgresql.service neo4j.service

[Install]
WantedBy=multi-user.target
